# systat.timer
[Unit]
Description=SyStat Matrix Timer

[Timer]
OnBootSec=60sec
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
