# systat.timer
[Unit]
Description=SyStat Matrix Timer

[Timer]
OnBootSec=10sec
OnUnitActiveSec=30sec
Persistent=true

[Install]
WantedBy=timers.target
